# custom-configs/nvidia-mobile-powermode/crontab
@reboot tmux new-session -d -s "nvidia-power" ~/bin/nvidia-powermode.sh

